# root/etc/services.d/nordvpn/finish
#!/bin/bash

RUNDIR="/run/nordvpn"
# IFACE="nordlynx"

rm -rf "${RUNDIR:?}"
# ip link delete "${IFACE:?}" 2>/dev/null
